# gateway/systemd/agentcart.service.example
[Unit]
Description=AgentCart gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=agentcart
Group=agentcart
WorkingDirectory=/opt/agentcart/gateway
Environment=HOME=/var/lib/agentcart
EnvironmentFile=/etc/agentcart/gateway.env
ExecStart=/usr/bin/python3 /opt/agentcart/gateway/agentcart.py
Restart=on-failure
RestartSec=3
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
