# gateway/systemd/agentcart-mpp-smoke.service.example
[Unit]
Description=AgentCart Tempo MPP smoke paid endpoint
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=openclaw
Group=openclaw
WorkingDirectory=/home/openclaw/workspace/agentcart
EnvironmentFile=-/etc/openclaw/agentcart-mpp-smoke.env
ExecStart=/usr/bin/node /home/openclaw/workspace/agentcart/scripts/mpp-smoke-server.mjs
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target

