# gateway/systemd/agentcart-mpp-smoke.service.example
[Unit]
Description=AgentCart Tempo MPP smoke paid endpoint
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=agentcart
Group=agentcart
WorkingDirectory=/opt/agentcart/gateway
EnvironmentFile=-/etc/agentcart/mpp-smoke.env
ExecStart=/usr/bin/node /opt/agentcart/gateway/scripts/mpp-smoke-server.mjs
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
